# roles/deploy_doca_firefly_container/templates/rc.local.j2
#!/bin/bash
ethtool --set-priv-flags {{ ptp_interface }} tx_port_ts on

